# collect/sentinel/sentinel-proxy/files/renew.cron
## crontab
#

# Random delay to spread checks/renewals across whole periods
RANDOM_DELAY=719

# Periodically check and renew Sentinel:Proxy certificate
0 */12 * * * root sentinel-certgen certs --hooks-dir /usr/libexec/sentinel/renew_hooks.d
