# net/morce/files/morce.cron
MAILTO=""
5 3 * * * root /etc/init.d/morce enabled && /etc/init.d/morce aggregate
