# hardware/omnia-ng/omnia-ng-support/files/wifi-checker.init
#!/bin/sh /etc/rc.common
START=99
STOP=01
USE_PROCD=1

start_service() {
    if uci show wireless | grep -q EHT[0-9]; then
        create_notification -s error 'Your Wi-Fi is set to Wi-Fi 7. We are currently experiencing issues with the Wi-Fi 7 driver. If you encounter instabilities, try temporarily downgrading to Wi-Fi 5.'
    fi
}

